# .codeploy/portable.fs/usr/lib/systemd/system/ops-redis.service
[Unit]
Wants=redis@ops\x2dredis.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c "echo provision ops-redis Redis server"
